# playbooks/roles/zsh/files/profile.d/alias_freebsd
alias ls='ls -G'
alias less='less -FXR'
