# playbooks/roles/zsh/files/profile.d/env
export EDITOR="vim"
export PAGER="less -FXR"
