# playbooks/roles/zsh/files/profile.d/env
export EDITOR="vim"
export PAGER="less -FXR"

# language settings
export LANG="en_US.UTF-8"
export LC_ALL="en_US.UTF-8"
export LC_MEASUREMENT="sv_SE.UTF-8"
export LC_TIME="sv_SE.UTF-8"
export LC_PAPER="sv_SE.UTF-8"
